# doc/gostaticserve.service
[Unit]
Description=put-file - High Performance Static File Server
After=network.target
Wants=network.target

[Service]
Type=simple

# Service running user
User=www-data
# Service running group
Group=www-data

# Working directory
WorkingDirectory=/opt/put-file

# Service start command
ExecStart=/opt/put-file/put-file \n--port=8080 \n--root=/opt/put-file/uploads \n--username=admin \n--password=your_secure_password \n--log-level=INFO \n--log-file=/var/log/put-file.log

# Restart policy
Restart=on-failure
RestartSec=5s

# Environment variables
Environment=GO_ENV=production

# Resource limits
LimitNOFILE=65535

# Log configuration
StandardOutput=append:/var/log/put-file.stdout.log
StandardError=append:/var/log/put-file.stderr.log

[Install]
WantedBy=multi-user.target

# Notes:
# 1. Make sure to create the running user and group (e.g., www-data)
# 2. Make sure to set correct file and directory permissions
# 3. Replace your_secure_password with the actual admin password
# 4. You can adjust the port, root directory and other parameters as needed
# 5. If you want to use a configuration file, you can replace the parameters with --config=/path/to/config.json
# 6. Install the service: sudo cp put-file.service /etc/systemd/system/
# 7. Enable the service: sudo systemctl enable put-file
# 8. Start the service: sudo systemctl start put-file
# 9. Check status: sudo systemctl status put-file